# templates/apache/apache.service.in
[Unit]
Description=Aura Platform Service: apache (@AURA_USER_NAME@)
After=network.target remote-fs.target nss-lookup.target
 
[Service]
Type=forking
PIDFile=/var/run/aura-@AURA_USER_CLASS@-@AURA_USER_NAME@-apache.pid
ExecStart=/usr/sbin/httpd -f /aura/configuration/@AURA_USER_CLASS@/@AURA_USER_NAME@/apache/apache.conf
ExecReload=/usr/sbin/httpd -k restart
ExecStop=/usr/sbin/httpd -k stop
PrivateTmp=true
 
[Install]
WantedBy=multi-user.target
